# etc/supervisor/conf.d/rfm2mqtt.conf
[program:rfm2mqtt]
command=rfm2mqtt.py
user=pi
autostart=true
autorestart=true
environment=PYTHONUNBUFFERED=1
